# skwikis.ini
[program:skwikis]
directory = /root/ShanghaiKingsWikis
command=gunicorn httpservice:app -c ./config.py
autostart=true
autorestart=true
user = root
stdout_logfile=/root/ShanghaiKingsWikis/logfile.log
std_logfile_maxbytes=50MB
std_logfile_backups=20