# skwikis.ini
[program:skwikis]
directory = /root/ShanghaiKingsWikis
command=gunicorn httpservice:app -c ./config.py
autostart=true
autorestart=true
user = root
logfile = /root/ShanghaiKingsWikis/logfile.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info